# homebrewupdate.plist
<!--  
  1. Place homebrewupdate.sh in /opt/nuna/scripts/homebrewupdate.sh
  2. Place the .plist file in ~/Library/LaunchAgents
  3. RUN: launchctl load ~/Library/LaunchAgents/homebrewupdate.plist
  NOTE: You might need to chmod +x homebrewupdate.sh
-->
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>homebrewupdate.plist</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/nuna/scripts/homebrewupdate.sh</string>
  </array>

  <key>Nice</key>
  <integer>1</integer>

  <key>StartInterval</key>
  <integer>21600</integer> <!-- UPDATE interval in seconds -->

  <key>RunAtLoad</key>
  <true/>

  <key>StandardErrorPath</key>
  <string>/tmp/nuna/homebrewupdate.err</string>

  <key>StandardOutPath</key>
  <string>/tmp/nuna/homebrewupdate.out</string>
</dict>
</plist>
